# desktop-agent/build-pkg/scripts/postinstall
#!/bin/bash
# .pkg 설치 직후 — 콘솔 사용자에게 LaunchAgent plist 를 생성하고 부트스트랩한다.
#
# 핵심: macOS Installer 는 root 로 postinstall 을 실행하므로
#       $USER / $HOME 이 root 를 가리킨다. 실제 사용자는 콘솔에서 가져와야 한다.

set -e

CONSOLE_USER=$(stat -f "%Su" /dev/console 2>/dev/null || echo "")
if [[ -z "$CONSOLE_USER" || "$CONSOLE_USER" == "root" || "$CONSOLE_USER" == "_mbsetupuser" ]]; then
    echo "⚠ 콘솔 사용자를 식별할 수 없습니다. 설치 후 사용자가 로그인한 상태에서:"
    echo "    cp /usr/local/share/aidesk/templates/com.aidesk.agent.plist ~/Library/LaunchAgents/"
    echo "    launchctl load ~/Library/LaunchAgents/com.aidesk.agent.plist"
    exit 0
fi

USER_UID=$(id -u "$CONSOLE_USER")
USER_HOME=$(eval echo "~$CONSOLE_USER")
LA_DIR="$USER_HOME/Library/LaunchAgents"
LOG_DIR="$USER_HOME/Library/Logs"
PLIST="$LA_DIR/com.aidesk.agent.plist"

# LANG 추출 — UTF-8 wide-char 처리에 필수 (한글 등).
# 우선순위: 1) 사용자 zsh 인터랙티브 (.zshrc 의 export 까지 거침)
#           2) macOS AppleLocale (UI 언어 설정) + ".UTF-8"
#           3) en_US.UTF-8 (최종 폴백)
USER_LANG=$(sudo -u "$CONSOLE_USER" -i zsh -ic 'echo "${LANG}"' 2>/dev/null | tail -1 | tr -d '[:space:]' || true)

if [[ -z "$USER_LANG" || "$USER_LANG" == "C" || "$USER_LANG" == "C.UTF-8" || "$USER_LANG" == "POSIX" ]]; then
    APPLE_LOCALE=$(sudo -u "$CONSOLE_USER" defaults read -g AppleLocale 2>/dev/null | tr -d '[:space:]' || true)
    [[ -n "$APPLE_LOCALE" ]] && USER_LANG="${APPLE_LOCALE}.UTF-8"
fi

[[ -z "$USER_LANG" ]] && USER_LANG="en_US.UTF-8"

sudo -u "$CONSOLE_USER" mkdir -p "$LA_DIR" "$LOG_DIR" "$USER_HOME/.aidesk"

# PyInstaller 바이너리를 직접 호출하는 plist 작성.
# (기존 uv-based plist 는 사용 안 함 — 바이너리 자체가 의존성을 안고 있음)
cat > "$PLIST" <<EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aidesk.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/aidesk-helper</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>ProcessType</key><string>Background</string>
  <key>StandardOutPath</key>
  <string>$LOG_DIR/aidesk-agent.log</string>
  <key>StandardErrorPath</key>
  <string>$LOG_DIR/aidesk-agent.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:$USER_HOME/.local/bin</string>
    <key>HOME</key>
    <string>$USER_HOME</string>
    <key>LANG</key>
    <string>$USER_LANG</string>
    <key>LC_CTYPE</key>
    <string>UTF-8</string>
  </dict>
</dict>
</plist>
EOF

chown "$CONSOLE_USER" "$PLIST"
chmod 644 "$PLIST"

# 부트스트랩 (사용자 GUI 도메인). 이미 등록되어 있으면 먼저 bootout.
launchctl bootout "gui/$USER_UID/com.aidesk.agent" 2>/dev/null || true
launchctl bootstrap "gui/$USER_UID" "$PLIST"

# ~/.claude.json 의 mcpServers 에 aidesk-channel 자동 등록 (멱등).
# 기존 aidesk-channel 등록이 있으면 binary 경로만 갱신, env 는 보존.
# 없으면 default sentinel UUID 로 새로 추가. claude code 재시작 후 적용.
AIDESK_CHANNEL_BIN="/usr/local/share/aidesk/aidesk-channel/bin/aidesk-channel"
CLAUDE_JSON="$USER_HOME/.claude.json"
if [[ -x "$AIDESK_CHANNEL_BIN" ]]; then
    sudo -u "$CONSOLE_USER" python3 - <<PYEOF
import json, os
path = os.path.expanduser("~/.claude.json")
try:
    with open(path) as f:
        data = json.load(f)
except FileNotFoundError:
    data = {}
except json.JSONDecodeError as e:
    print(f"⚠ ~/.claude.json 파싱 실패 ({e}) — aidesk-channel 자동 등록 건너뜀")
    raise SystemExit(0)

servers = data.setdefault("mcpServers", {})
existing = servers.get("aidesk-channel", {})
existing_env = (existing.get("env") or {}) if isinstance(existing, dict) else {}

servers["aidesk-channel"] = {
    "command": "node",
    "args": ["$AIDESK_CHANNEL_BIN"],
    "env": {
        # 기존 env 보존, 없는 키만 default 채움.
        "AIDESK_AGENT_ID": existing_env.get("AIDESK_AGENT_ID", "a1b2c3d4-0000-0000-0000-000000000001"),
        "AIDESK_API_URL":  existing_env.get("AIDESK_API_URL", "http://localhost:30081"),
    },
}
with open(path, "w") as f:
    json.dump(data, f, indent=2, ensure_ascii=False)
    f.write("\n")
print("✓ ~/.claude.json 에 aidesk-channel mcp 등록 (claude code 재시작 후 적용)")
PYEOF
else
    echo "⚠ aidesk-channel binary 없음 ($AIDESK_CHANNEL_BIN) — mcp 자동 등록 건너뜀"
fi

echo "✓ AI Desk Helper 설치 완료"
echo "  - 바이너리: /usr/local/bin/aidesk-helper"
echo "  - LaunchAgent: $PLIST"
echo "  - 로그: $LOG_DIR/aidesk-agent.err"
echo "  - aidesk-channel mcp: $AIDESK_CHANNEL_BIN"

exit 0
